# dnsmasq.ini
[program:dnsmasq]
user=root
command=/usr/sbin/dnsmasq -k
