# pihole.ini
[program:pihole]
user=pihole
command=/usr/bin/pihole-FTL
